# mark-ii-raspberrypi/files/etc/systemd/system/mycroft-hal.service
[Unit]
Description=DBus server for hardware access

[Service]
Type=simple
ExecStart=/opt/mycroft/dbus-hal/run.sh

[Install]
WantedBy=default.target
