# docker/files/etc/systemd/system/mycroft-switch.service
[Unit]
Description=Run switch service
After=mycroft-xmos.service

[Service]
User=pi
WorkingDirectory=/home/mycroft
ExecStart=/usr/bin/mark2-switch-server

[Install]
WantedBy=default.target
