# docker/files/etc/systemd/system/mycroft-volume.service
[Unit]
Description=Run switch volume service
After=mycroft-switch.service

[Service]
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/mark2-volume-client
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
