# docker/files/etc/systemd/system/mycroft-switch.service
[Unit]
Description=Run switch service
After=mycroft-xmos.service

[Service]
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/mark2-switch-server
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
