# docker/files/etc/systemd/system/mycroft-hardware-dbus.service
[Unit]
Description=DBus server for hardware access

[Service]
ExecStart=/opt/mycroft-mark2-hardware/run.sh

[Install]
WantedBy=default.target
